# archive-hardlink.timer
[Unit]
Description=Monthly Archive Hardlinkng

[Timer]
OnCalendar=*-*-03 03:00:00
AccuracySec=1m
Persistent=true

[Install]
WantedBy=timers.target
