# archive.service
[Unit]
Description=Archive Update

[Service]
Type=oneshot
User=archive
Group=archive
ExecStart=/usr/bin/archive
