# archive-repos.timer
[Unit]
Description=Daily Archlinux Archive Update

[Timer]
OnCalendar=12:42:00
AccuracySec=1m
Persistent=true

[Install]
WantedBy=timers.target
