# archive-aur.service
[Unit]
Description=Archive AUR

[Service]
Type=oneshot
ExecStart=/home/archlinux/archivetools/aur-rsync
