# archive-iso.service
[Unit]
Description=Archlinux Archive ISO Updates

[Service]
Type=oneshot
ExecStart=/home/archlinux/archivetools/iso-rsync
