# arm.service
[Unit]
Description=ARM Update
After=archive-repos.service

[Service]
Type=oneshot
ExecStart=/home/archlinux/archivetools/arm-link
