# systemd/arch-mirror-archive.timer
[Unit]
Description=Daily Arch Repository Sync

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=multi-user.target
