# systemd/arch-mirror-archive.service
[Unit]
Description=Arch Mirror Synchronization Service

[Service]
Type=oneshot
ExecStart=/opt/local/archivetools/bin/sync
